# systemd/tokenchaind-testnet.service
[Unit]
Description=TokenChain Testnet Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tokenchain
Group=tokenchain
WorkingDirectory=/var/lib/tokenchain-testnet
EnvironmentFile=-/etc/tokenchain/tokenchaind-testnet.env
ExecStart=/usr/local/bin/tokenchaind start --home /var/lib/tokenchain-testnet
Restart=on-failure
RestartSec=3
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/tokenchain-testnet

[Install]
WantedBy=multi-user.target
