# systemd/tokenchaind.service
[Unit]
Description=TokenChain Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tokenchain
Group=tokenchain
WorkingDirectory=/var/lib/tokenchain-mainnet
EnvironmentFile=-/etc/tokenchain/tokenchaind-mainnet.env
ExecStart=/usr/local/bin/tokenchaind start --home /var/lib/tokenchain-mainnet
Restart=on-failure
RestartSec=3
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/tokenchain-mainnet

[Install]
WantedBy=multi-user.target
