# systemd/tokenchain-indexer.service
[Unit]
Description=TokenChain Indexer API
After=network-online.target tokenchaind-testnet.service
Wants=network-online.target

[Service]
Type=simple
User=tony
Group=tony
WorkingDirectory=/var/www/tokenchain-indexer
ExecStart=/usr/local/bin/tokenchain-indexer
Restart=always
RestartSec=2
Environment=LISTEN_ADDR=127.0.0.1:3321
Environment=CHAIN_ID=tokenchain-testnet-1
Environment=NETWORK=testnet
Environment=RPC_ADDR=http://127.0.0.1:26657
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/tokenchain-indexer

[Install]
WantedBy=multi-user.target
