# systemd/tokenchain-relayer.service
[Unit]
Description=TokenChain Hermes Relayer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tokenchain
WorkingDirectory=/opt/tokenchain-ops/relayer
ExecStart=/usr/local/bin/hermes start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
